# scripts/trading-worker.service
[Unit]
Description=Trading Platform Worker (24/7 scheduler)
After=network.target ibgateway-live.service
Wants=ibgateway-live.service
# StartLimit* keys must live in [Unit] section, not [Service] (systemd v229+)
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
User=root
WorkingDirectory=/opt/trading-platform
ExecStart=/opt/trading-platform/.venv/bin/python worker.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/trading-platform/.env

# Logging
StandardOutput=append:/opt/trading-platform/logs/worker/worker_systemd.log
StandardError=append:/opt/trading-platform/logs/worker/worker_systemd.log

# Safety
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
